# crontab
*/1 * * * * root /usr/bin/node /etc/healthcheck.js
